# reminders.service
[Unit]
Description=Reminders Notifier
Wants=default.target

[Service]
Type=simple
ExecStart=reminders
